# crontab
*/2 * * * * root /root/Lab2/lab2-part3.sh
